# etc/systemd/system/isutar.go.service
[Unit]
Description = isucon6 qualifier sub application in Go

[Service]
WorkingDirectory=/home/isucon/deploy/go
EnvironmentFile=/home/isucon/env.sh

ExecStart = /home/isucon/deploy/go/isutar

Restart   = always
Type      = simple
User      = isucon
Group     = isucon

[Install]
WantedBy = multi-user.target
